# scripts/paxi_supervisor.conf
[program:paxi-uvicorn]
directory=/root/utils_server
command=/root/utils_server/venv/bin/uvicorn utils_server:app --host 127.0.0.1 --port 8080
autostart=true
autorestart=true
user=root
